# power_manager/udev/70-wifi-powersave.rules
# /etc/udev/rules.d/70-wifi-powersave.rules

ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/iw dev %k set power_save on"