# power_manager/udev/71-hdmi_sound.rules
ACTION=="change", SUBSYSTEM=="drm", KERNEL=="card*", RUN+="/usr/bin/echo `date` >> /tmp/drm.log"
