# power_manager/udev/71-hdmi_sound.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/hdmi_sound_toggle"
